# supervisord.conf
[supervisord]
nodaemon=true

[program:jenkins]
directory=/usr/share/jenkins
command="java", "-jar", "jenkins.war"
stopasgroup=true
killasgroup=true
stderr_logfile=/var/log/jenkins/stderr.log
